# files/etc/init.d/languard-pro
#!/bin/sh /etc/rc.common

START=98
STOP=10

start() {
    echo "Starting LanGuard Pro"
    [ -x /root/languard-sync-devices.sh ] && /root/languard-sync-devices.sh >/dev/null 2>&1
    [ -x /root/languard-apply-blocks.sh ] && /root/languard-apply-blocks.sh >/dev/null 2>&1
    [ -x /root/languard-limits-engine.sh ] && /root/languard-limits-engine.sh >/dev/null 2>&1
}

stop() {
    echo "Stopping LanGuard Pro"
}

restart() {
    start
}
